# packages/ambient-edge-server/ambient_edge_server-0.3.0.tar.gz/ambient_edge_server-0.3.0/ambient_edge_server/templates/ambient_edge_server.service.jinja2
[Unit]
Description={{ service.description }}
After=network.target

[Service]
Type=simple
User={{ service.user }}
Group={{ service.group }}
ExecStart={{ service.exec_start }}
Restart=on-failure
{% if service.environment %}
Environment={{ service.environment }}
{% endif %}
# WorkingDirectory={{ service.working_directory }}
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target